# programs/sophiatx_light/package/debian/postinst
#!/bin/sh
# postinst script for sophiatx-light
#
# see: dh_installdeb(1)

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
# <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
# <failed-install-package> <version> `removing'
# <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
        DATA_DIR=/var/lib/sophiatx-light
        LOG_DIR=/var/log/sophiatx-light
        LOG_FILE=$LOG_DIR/sophiatx-light.log

        command -v systemctl >/dev/null 2>&1
        if [ "$?" -ne 0 ]; then
            echo "System must support systemctl to run sophiatx-light." >&2
            exit 1
        fi

        id sophiatx-light >/dev/null 2>&1
        if [ "$?" -ne 0 ]; then
            useradd --system -U -M sophiatx-light -s /bin/false -d $DATA_DIR
        fi

        test -d $DATA_DIR || mkdir -p $DATA_DIR
        test -d $LOG_DIR || mkdir -p $LOG_DIR
        touch $LOG_FILE

        chown -R -L sophiatx-light:sophiatx-light $DATA_DIR
        chown -R -L syslog:adm $LOG_DIR
        chmod 755 $DATA_DIR
        chmod 755 $LOG_DIR
        chmod 644 $LOG_FILE

        # Add defaults file, if it doesn't exist
        if [ ! -e /etc/default/sophiatx-light ]; then
            touch /etc/default/sophiatx-light
        fi

        systemctl enable sophiatx-light
     ;;

     abort-upgrade|abort-remove|abort-deconfigure)
     ;;

     *)
         echo "postinst called with unknown argument \`$1'" >&2
         exit 1
     ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0