# programs/sophiatxd/package/debian/postinst
#!/bin/sh
# postinst script for sophiatx
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
# <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
# <failed-install-package> <version> `removing'
# <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
        DATA_DIR=/var/lib/sophiatx
        LOG_DIR=/var/log/sophiatx
        SCRIPT_DIR=/usr/lib/sophiatx/scripts

        command -v systemctl >/dev/null 2>&1
        if [ "$?" -ne 0 ]; then
            echo "System must support systemctl to run sophiatx." >&2
            exit 1
        fi

        id sophiatx >/dev/null 2>&1
        if [ "$?" -ne 0 ]; then
            useradd --system -U -M sophiatx -s /bin/false -d $DATA_DIR
        fi



        test -d $DATA_DIR || mkdir -p $DATA_DIR
        test -d $LOG_DIR || mkdir -p $LOG_DIR

        chown -R -L sophiatx:sophiatx $DATA_DIR
        chown -R -L sophiatx:sophiatx $LOG_DIR
        chmod 755 $DATA_DIR
        chmod 755 $LOG_DIR

        # Add defaults file, if it doesn't exist
        if [[ ! -f /etc/default/sophiatx ]]; then
            touch /etc/default/sophiatx
        fi

        cp -f $SCRIPT_DIR/sophiatx.service /lib/systemd/system/sophiatx.service
        systemctl enable sophiatx
     ;;

     abort-upgrade|abort-remove|abort-deconfigure)
     ;;

     *)
         echo "postinst called with unknown argument \`$1'" >&2
         exit 1
     ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0